# com.peter-brandt.munger.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.peter-brandt.munger</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>$HOME/Documents/claudcode/math-curve-loaders/trigger_update.sh munger</string>
    </array>
    <!-- Saturday 09:00 local (weekly scan after Friday close) -->
    <key>StartCalendarInterval</key>
    <dict><key>Weekday</key><integer>6</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <key>StandardOutPath</key><string>/tmp/peter-brandt-munger.out</string>
    <key>StandardErrorPath</key><string>/tmp/peter-brandt-munger.err</string>
    <key>EnvironmentVariables</key>
    <dict><key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string></dict>
</dict>
</plist>
