# com.peter-brandt.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.peter-brandt.daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>$HOME/Documents/claudcode/math-curve-loaders/trigger_update.sh daily</string>
    </array>
    <!-- Tue–Sat 07:00 local (≈1h after US close, yfinance has prev-day bar) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>7</integer></dict>
        <dict><key>Weekday</key><integer>6</integer><key>Hour</key><integer>7</integer></dict>
    </array>
    <key>StandardOutPath</key><string>/tmp/peter-brandt-daily.out</string>
    <key>StandardErrorPath</key><string>/tmp/peter-brandt-daily.err</string>
    <key>EnvironmentVariables</key>
    <dict><key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string></dict>
</dict>
</plist>
